# system_files/usr/lib/systemd/system/dual-wan.timer
[Unit]
Description=Verifica dual WAN a cada 30s

[Timer]
OnBootSec=15
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
